# modbus-bridge.service
# Systemd service unit — reference copy
# The installer (install.sh) writes the live version to /etc/systemd/system/
# with the actual User, WorkingDirectory, and Node.js path substituted in.
#
# To install manually:
#   sudo cp modbus-bridge.service /etc/systemd/system/
#   # edit User= and WorkingDirectory= to match your setup
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now modbus-bridge

[Unit]
Description=Modbus Bridge
Documentation=https://github.com/matnig/direktvermarktung-virtuallogger-modbustcp-vpn
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/modbus-bridge
ExecStart=/usr/bin/node src/app/main.js
Restart=on-failure
RestartSec=5

Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=3000

# Forward stdout/stderr to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=modbus-bridge

[Install]
WantedBy=multi-user.target
